# new-api.service
# File path: /etc/systemd/system/new-api.service
# sudo systemctl daemon-reload
# sudo systemctl start new-api
# sudo systemctl enable new-api
# sudo systemctl status new-api
[Unit]
Description=Yihi API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/www/wwwroot/yihi-api
ExecStart=/www/wwwroot/yihi-api/yihi-api --port 3000 --log-dir /www/wwwroot/yihi-api/logs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target